# startup/.config/autostart/autostart.desktop
[Desktop Entry]
Name=Startup script
Type=Application
Exec=/home/jisoo/bin/autostart.sh
Terminal=false
X-GNOME-Autostart-enabled=true
